# prod/twitchbot.service
[Unit]
Description=BotName

[Service]
Type=simple
ExecStart=/<PATH TO FOLDER>/.venv/bin/python /<PATH TO FOLDER>/bot/bot_application.py

[Install]
WantedBy=multi-user.target
